# script/oi-config-server.service
[Unit]
Description=One Invest Config Server

[Service]
ExecStart=/home/temp/oi_config.sh

[Install]
WantedBy=multi-user.target
