# script/oi-main-logic-server.service
[Unit]
Description=One Invest Main Logic Server

[Service]
ExecStart=/home/temp/oi_main_logic.sh

[Install]
WantedBy=multi-user.target
