# script/oi-notification-server.service
[Unit]
Description=One Invest Notification Server

[Service]
ExecStart=/home/temp/oi_notification.sh

[Install]
WantedBy=multi-user.target
