# script/oi-notification-server-8763.service
[Unit]
Description=One Invest Notification Server 8763

[Service]
ExecStart=/home/temp/oi_notification_8763.sh

[Install]
WantedBy=multi-user.target
